# init/poi-daemon.ini
[program:poi-daemon]
command=/usr/local/bin/poi -c /path/to/config daemon
autostart=true
stopsignal=KILL
