# init/poi_daemon.service
[Unit]
Description=Poi daemon

[Service]
ExecStart=/usr/local/bin/poi -c /path/to/config daemon
Restart=always

[Install]
WantedBy=multi-user.target
